# salt-master.conf
[program:salt-master]
directory=/
command=/usr/bin/salt-master -d
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log